# debian/ve-drv-kmod.postinst
#!/bin/bash
set -e
depmod -a

#DEBHELPER#

/bin/systemctl enable ve_drv > /dev/null 2>&1

exit 0
